# systemd/rcon-webgui.service
[Unit]
Description=Minecraft WebGUI (RCON/log panel)
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/mc-rcon-webgui/app
EnvironmentFile=/opt/mc-rcon-webgui/app/.env
ExecStart=/usr/bin/node /opt/mc-rcon-webgui/app/server.js
Restart=always
RestartSec=2
# NOTE: runs as root by default (matches your helper scripts).
# For hardening, set User=mcgui and grant sudo below.

[Install]
WantedBy=multi-user.target
